# universal-telemetry-software/deploy/car-telemetry.service
# DAQ Telemetry — Car RPi native systemd service
#
# Runs the Python telemetry stack directly on the car Pi (no Docker, no Redis).
# Redis is optional — its only role on the car is serving the local WebSocket
# bridge so PECAN can connect directly to the car. That's never needed in
# normal operation. Without Redis the car still reads CAN, sends UDP, logs
# Parquet, and serves the status page on port 8080.
#
# Install:
#   sudo cp car-telemetry.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable car-telemetry
#   sudo systemctl start car-telemetry

[Unit]
Description=DAQ Telemetry Car
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/daq-radio/universal-telemetry-software

# Hardcoded car environment — edit here if IPs change
Environment=ROLE=car
Environment=REMOTE_IP=10.71.1.20
Environment=UDP_PORT=5005
Environment=TCP_PORT=5006
Environment=REDIS_URL=redis://127.0.0.1:6379/0
Environment=WS_PORT=9080
Environment=STATUS_PORT=8080
Environment=SIMULATE=false
Environment=ENABLE_VIDEO=false
Environment=ENABLE_AUDIO=false
Environment=SET_TIME_ENABLED=true
Environment=SHUTDOWN_ENABLED=true
Environment=ENABLE_INFLUX_LOGGING=false
# Set at install time: sed -i "s/GIT_HASH=unknown/GIT_HASH=$(git rev-parse --short HEAD)/" car-telemetry.service
Environment=GIT_HASH=unknown
# No REDIS_URL — Redis is not needed on the car. Without it, publishes no-op
# and the local WebSocket bridge won't serve data, but UDP streaming to the
# base station is unaffected.

# uv run handles virtualenv creation and dependency sync automatically.
# Run `uv sync` once during setup (requires internet) to pre-build the venv.
ExecStart=/home/pi/.local/bin/uv run python main.py

# Always restart — handles CAN bus glitches, Redis not ready yet, etc.
Restart=always
RestartSec=3

# Give the process CAP_SYS_NICE so it can set its own scheduler priority
AmbientCapabilities=CAP_SYS_NICE
CapabilityBoundingSet=CAP_SYS_NICE

# Log to systemd journal — view with: journalctl -u car-telemetry -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=car-telemetry

[Install]
WantedBy=multi-user.target
